# schedules.cron
20,40,59 * * * * echo 3 > /proc/sys/vm/drop_caches
0 0 * * * sh /root/gospely/deploy/dbback.sh
0 2 * * * sh /root/gospely/deploy/vote_back.sh
0 3 * * * curl -s "http://localhost:9999?token=gospelytokenaasss&operation=deleteDocker"
20,40,59 * * * * curl -s "http://localhost:9999/schedules?token=gospelytokenaasss&operation=clearInnersession"
10,20,30,40,50,59 * * * * docker rm -f $( docker stop $(docker ps -a | grep -v 'gospel' | awk '{print $1}'))
1  3  * * * freshclam
30 4  * * * clamscan  --infected  -r /  --remove -l /var/log/clamscan.log
