# schedules.cron
20,40,59 * * * * echo 3 > /proc/sys/vm/drop_caches
