# schedules.cron
20,40,59 * * * * echo 3 > /proc/sys/vm/drop_caches
00,00,01 * * * sh /root/gospely/deploy/dbback.sh
01,00,01, * * * sh /root/gospely/deploy/vote_back.sh
01,00,30, * * * curl http://localhost:9999?token=gospelytokenaasss
