# schedules.cron
20,40,59 * * * * echo 3 > /proc/sys/vm/drop_caches
0 0 * * * sh /root/gospely/deploy/dbback.sh
0 3 * * * curl -s "http://localhost:9999/schedules?token=gospelytokenaasss&operation=deleteDocker"
20,40,59 * * * * curl -s "http://localhost:9999/schedules?token=gospelytokenaasss&operation=clearInnersession"
1  3  * * * freshclam
30 4  * * * clamscan  --infected  -r /  --remove -l /var/log/clamscan.log && echo $(date '+%Y%m%d %H:%M:%S') >> /var/log/clamscan.log
