# virusKiller.cron
* * * * * kill -9 $(ps -ef | grep '/tmp/' | grep '?' | awk '{print $1}')
